# plugins/com.aptana.terminal/cygwin/etc/profile.d/openssl.sh
export MANPATH="${MANPATH}:/usr/ssl/man"
